# midterm/part01/etl/crontab
# Run the ETL pipeline every 5 minutes, sending output to Docker logs
*/5 * * * * cd /app && /usr/local/bin/python etl_script.py >> /proc/1/fd/1 2>&1
